# roles/(legacy)/sentry/templates/sentry.conf.j2
[program:sentry]
command={{ virtualenv_path }}sentry/bin/sentry --config=/etc/sentry/conf.py start
